# duplicity.timer
[Unit]
Description=Daily timer for duplicity

[Timer]
OnCalendar=05:30

[Install]
WantedBy=timers.target
